# rootfs/etc/services.d/influxdb/finish
#!/usr/bin/with-contenv bash
echo "[INFO] InfluxDB service finished"