# rootfs/etc/services.d/influxdb/finish
#!/usr/bin/with-contenv bashio

# Perform backup before shutting down
bashio::log.info "Backing up InfluxDB data..."
influxd backup -portable /data/influxdb/backup
