# rootfs/etc/services.d/influxdb/finish
#!/usr/bin/with-contenv bashio

# Log the exit reason
bashio::log.warning "InfluxDB service exited. Attempting restart..."

# Optional: Add any cleanup or pre-restart tasks here
exit 0