# rootfs/etc/services.d/influxdb/run
#!/usr/bin/with-contenv bashio

# Create directories with proper permissions if they don't exist
mkdir -p /data/influxdb/meta /data/influxdb/data /data/influxdb/wal
chown -R nobody:nobody /data/influxdb
chmod -R 777 /data/influxdb

# Start InfluxDB with config
exec s6-setuidgid nobody influxd -config /etc/influxdb/influxdb.conf