# rootfs/etc/services.d/carbonoz/finish
#!/usr/bin/with-contenv bash
echo "[INFO] Carbonoz service finished"